# etc/flower.socket
[Unit]
Description=Flower socket

[Socket]
ListenStream=0.0.0.0:5555

[Install]
WantedBy=sockets.target
